# example/kubernetes/docker/ml/pvacms-cluster-supervisor.conf
[program:pvacms]

; ---------------------------------------------------------------------------
; Environment:
;   KRB5_KTNAME / KRB5_CLIENT_KTNAME    Kerberos keytab used by pvacms to
;                                        obtain service credentials for the
;                                        pvacms/cluster@EPICS.ORG principal.
;   EPICS_AUTH_KRB_REALM                Default realm assumed for unqualified
;                                        principals in Kerberos auth flows.
;                                        Assumes network allows access to KDC.
;   EPICS_PVACMS_CLUSTER_NAME_SERVERS   Cross-zone PVACMS cluster replication
;                                        gateway: points at the LAB zone's
;                                        gateway so this ML-zone PVACMS can
;                                        reach the lab PVACMS for cluster SYNC.
; ---------------------------------------------------------------------------
environment=
    KRB5_KTNAME="/home/idm/.config/krb5/pvacms.keytab",
    KRB5_CLIENT_KTNAME="/home/idm/.config/krb5/pvacms.keytab",
    EPICS_AUTH_KRB_REALM="EPICS.ORG",
    EPICS_PVACMS_CLUSTER_NAME_SERVERS="pvxs-lab-gateway-xgw:5175"

; ---------------------------------------------------------------------------
; Command:
;   --cluster-mode                      Run PVACMS as a cluster node (enables
;                                        peer discovery, SYNC, leader election).
;   --cluster-discovery-timeout 30      Wait 30s for cluster discovery because
;                                        discovery via gateway is slower than
;                                        normal (default 10s).
;   --cluster-bidi-timeout 30           Extend bidirectional-connectivity
;                                        probe during cluster join from the
;                                        default 5s to 30s; the probe opens
;                                        a test subscription back through the
;                                        gateway, which adds multi-hop latency.
;   --cert-auth-keychain <file>         Path to the CA keychain (.p12)
;                                        created in the ca-keygen job.
;   --acf <file>                        ACF file controlling who may perform
;                                        administrative cert operations.
;
;   NOTE: --admin-keychain-ensure is intentionally omitted; this ML-zone
;         PVACMS has no local admin user. Admin bootstrap happens on the
;         idm pod in the lab zone and replicates in via cluster SYNC.
;
;   `exec` replaces the shell PID so stopsignal=TERM reaches pvacms directly.
; ---------------------------------------------------------------------------
command=/bin/sh -c 'exec /opt/epics/pvxs-cms/bin/$(/opt/epics/epics-base/startup/EpicsHostArch)/pvacms --cluster-mode --cluster-discovery-timeout 30 --cluster-bidi-timeout 30 --cert-auth-keychain /etc/pvacms/cert_auth.p12 --acf /etc/pvacms/pvacms.acf -v'

user=idm
autostart=true
autorestart=true

; Send logs to container stdout/stderr
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true
